# scripts/newsletter/newsletter-send.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily newsletter send — job starts 5:50 AM PT so the email lands ~6:00.
  (The Fable editorial pass takes ~7-8 minutes; a timeout degrades to the
  deterministic build, never a missed send.)

  scheduled-send.mjs holds the shared Mini repo lock, fetches origin/main,
  safely fast-forwards a clean behind checkout, and fails closed on dirty,
  diverged, or locally modified source before it invokes send.mjs.

  Install on the Mini:
    cp scripts/newsletter/newsletter-send.plist ~/Library/LaunchAgents/org.southbaytoday.newsletter-send.plist
    launchctl bootout gui/$(id -u)/org.southbaytoday.newsletter-send 2>/dev/null || true
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/org.southbaytoday.newsletter-send.plist

  Logs: /tmp/sbt-newsletter-send.{out,err}.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.southbaytoday.newsletter-send</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/stephenstanwood/Projects/southbaytoday.org/scripts/newsletter/scheduled-send.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/stephenstanwood/Projects/southbaytoday.org</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>50</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/sbt-newsletter-send.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/sbt-newsletter-send.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
