# scripts/social/nightly-purge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.southbaytoday.nightly-purge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--env-file=.env.local</string>
    <string>scripts/social/nuke-old-posts.mjs</string>
    <string>--all-before-now</string>
    <string>--notify</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/stephenstanwood/Projects/southbaytoday.org</string>
  <!-- Midnight PT every night. Every day is a new day (south bay today!) —
       the profiles get a clean slate before the 6:15am morning post goes out. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/sbt-nightly-purge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sbt-nightly-purge.log</string>
</dict>
</plist>
