# scripts/social/publish-jitter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.southbaytoday.publish-jitter</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--env-file=.env.local</string>
    <string>scripts/social/regenerate-publish-plist.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/stephenstanwood/Projects/southbaytoday.org</string>
  <!-- Fires daily at 5:00 AM PT — comfortably before the earliest possible
       day-plan fire (06:15) so today's randomized plist is in place. -->
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key>
  <string>/tmp/sbs-publish-jitter.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sbs-publish-jitter.log</string>
</dict>
</plist>
