# scripts/social/weekly-purge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.southbaytoday.weekly-purge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>--env-file=.env.local</string>
    <string>scripts/social/nuke-old-posts.mjs</string>
    <string>--all-before-now</string>
    <string>--notify</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/stephenstanwood/Projects/southbaytoday.org</string>
  <!-- Saturday 3:30am PT — well before Stephen wakes at 4am. Each Saturday
       wipes the previous week so the profiles look fresh going into the
       weekend. New posts published Saturday morning (6:15am+) stay. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/sbt-weekly-purge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sbt-weekly-purge.log</string>
</dict>
</plist>
